# files/etc/rc.local
#!/bin/sh

# 修复 R2C Plus (eth0) 网卡
ethtool -K eth0 generic-receive-offload off
ethtool -K eth0 tx-scatter-gather off

exit 0
